# systemd/x-agent-engagement.timer
[Unit]
Description=Daily vibe-x-agent engagement refresh
Requires=x-agent-engagement.service

[Timer]
# Once a day, well clear of the posting windows so the two never race on
# the state file. Persistent=true catches missed runs (e.g. a reboot).
OnCalendar=*-*-* 06:00:00
Persistent=true
Unit=x-agent-engagement.service

[Install]
WantedBy=timers.target
